# ops/upstart/consul_client.conf
description "Consul"

start on vagrant-ready or runlevel [2345]
stop on runlevel [!2345]

respawn

console log

script
  if [ -f "/etc/service/consul" ]; then
    . /etc/service/consul
  fi

  # Make sure to use all our CPUs, because Consul can block a scheduler thread
  export GOMAXPROCS=`nproc`

  exec /usr/bin/consul agent \
   -data-dir="/tmp/consul" \
   -config-dir="/etc/consul.d" \
   -atlas-token="YOUR_ATLAS_TOKEN" \
   -atlas="<username>/infrastructure" \
   -atlas-join \
   ${CONSUL_FLAGS} \
   >>/var/log/consul.log 2>&1
end script
